# automation_studio/rootfs/etc/services.d/automation_studio/run
#!/usr/bin/with-contenv bashio
set -e

exec /app/run.sh
