# automation_studio/rootfs/etc/services.d/automation_studio/run
#!/usr/bin/env bash
set -e

exec /app/run.sh
